# timers/journal.timer
[Unit]
Description=Launch a daily journal later in the evening

[Timer]
OnCalendar=*-*-* 20:00:00
Persistent=true

[Install]
WantedBy=timers.target
