# timers/template.timer
[Unit]
Description=

[Timer]
OnCalendar=
Persistent=

[Install]
WantedBy=timers.target
